# expo-build-service/go-server.service
[Unit]
Description=Go Server Service
After=network.target

[Service]
Type=simple
User=distro
Group=distro
WorkingDirectory=/home/distro/Go/expo-build-server
ExecStart=/home/distro/Go/expo-build-server/buildHandler
Restart=on-failure
RestartSec=5s
Environment=PORT=8080

[Install]
WantedBy=multi-user.target
